# systemd/do-dyndns.timer
[Unit]
Description=Run Digital Ocean Dynamic DNS updater

[Timer]
OnCalendar=*:0/15
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
